# business-crm.service
[Unit]
Description=Business CRM Backend API
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/business-crm
Environment="NODE_ENV=production"
Environment="PORT=5000"
ExecStart=/usr/bin/node backend/server.js
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=business-crm

[Install]
WantedBy=multi-user.target
